# booster-command-server/scripts/booster-command-server.service
[Unit]
Description=Booster Command Server
After=network.target

[Service]
Type=simple
User=booster
Group=booster
WorkingDirectory=/opt/booster-command-server
Environment="PATH=/opt/booster-command-server/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="ROBOT=mock"
ExecStart=/opt/booster-command-server/venv/bin/python -m uvicorn main:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=10
StandardOutput=append:/var/log/booster-command-server/server.log
StandardError=append:/var/log/booster-command-server/server.log

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true

[Install]
WantedBy=multi-user.target

